# systemd/netbird-update.service
[Unit]
Description=NetBird Self-Hosted Automatic Update
Documentation=https://github.com/tolakang/netbird-autoupdate
After=docker.service network-online.target
Requires=docker.service

[Service]
Type=oneshot
# Default installation path is /opt/netbird. To use a different directory,
# edit this file and change ExecStart, Environment, and ReadWritePaths.
ExecStart=/opt/netbird/scripts/update-netbird.sh
Environment=COMPOSE_DIR=/opt/netbird
Nice=10
IOSchedulingClass=best-effort
PrivateTmp=yes
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/opt/netbird /opt/netbird/backups /run
User=root
Group=root

[Install]
WantedBy=timers.target